# profile.d/less-improved.sh
export LESS="$LESS --raw-control-chars --ignore-case --quit-if-one-screen --LONG-PROMPT --shift=5"
export SYSTEMD_LESS=$LESS
